# deploy/escudo-gateway.service
[Unit]
Description=Escudo VPN Gateway
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/escudo-gateway --config /etc/escudo/gateway.toml
Restart=always
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
