# deploy/escudo-admin.service
[Unit]
Description=Escudo Admin API
After=network.target postgresql.service

[Service]
Type=simple
User=escudo
ExecStart=/usr/local/bin/escudo-admin --config /etc/escudo/api.toml
Restart=always
RestartSec=5
Environment=RUST_LOG=info
EnvironmentFile=-/etc/escudo/api.env

[Install]
WantedBy=multi-user.target
